# TheWordBatch/crontab
25 4 * * * java -jar /TheWordBatch-1.0.jar
